# bin/control
#!/bin/sh
# chkconfig:         2345 75 15
# description:       sfs2x-service
### BEGIN INIT INFO
# Provides:          sfs2x-service
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: sfs2x-service
### END INIT INFO

SFX_SERVICE="/opt/repo/versions/2.x/SmartFoxServer_2X/SFS2X/sfs2x-service";

sed 's/127.0.0.1/<%= ENV['J_OPENSHIFT_SFS_IP'] %>/g'  /opt/repo/versions/2.x/SmartFoxServer_2X/SFS2X/config/server.xml > /opt/repo/versions/2.x/SmartFoxServer_2X/SFS2X/server.xml.erb
erb  /opt/repo/versions/2.x/SmartFoxServer_2X/SFS2X/server.xml.erb > /opt/repo/versions/2.x/SmartFoxServer_2X/SFS2X/config/server.xml

$SFX_SERVICE $1;
